# packaging/systemd/sure-smartie-linux.service.in
[Unit]
Description=Sure Smartie Linux LCD service
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=@SURE_SMARTIE_SERVICE_USER@
Group=@SURE_SMARTIE_SERVICE_GROUP@
SupplementaryGroups=@SURE_SMARTIE_SERIAL_GROUP@
EnvironmentFile=-@SURE_SMARTIE_ENV_FILE@
StateDirectory=sure-smartie-linux
LogsDirectory=sure-smartie-linux
WorkingDirectory=/var/lib/sure-smartie-linux
StandardOutput=journal
StandardError=journal
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/sure-smartie-linux
Restart=on-failure
RestartSec=2
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
