# build-desktop-test/sure-smartie-linux-root.service
[Unit]
Description=Sure Smartie Linux LCD service (root for Intel RAPL access)
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
EnvironmentFile=-/home/re/src/Sure-Smartie-linux/build/install-desktop-test/etc/default/sure-smartie-linux
StateDirectory=sure-smartie-linux
LogsDirectory=sure-smartie-linux
WorkingDirectory=/var/lib/sure-smartie-linux
StandardOutput=journal
StandardError=journal
ExecStart=/home/re/src/Sure-Smartie-linux/build/install-desktop-test/bin/sure-smartie-linux
Restart=on-failure
RestartSec=2
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
